# services/hyprland-graphical-session.target
[Unit]
Description=Hyprland Graphical Session Trigger
BindTo=graphical-session.target

[Install]
WantedBy=graphical-session.target
